# python-script-host/rc.local
#!/bin/sh
/opt/splunk/bin/splunk start --accept-license
exit 0