# splunk/rc.local
#!/bin/sh
/opt/splunk/bin/splunk start
exit 0