# airvideohd/rc.local
#!/bin/sh
/AirVideoServerHD --config=/config/server.properties
exit 0